# prometheus/node_exporter.service
[Unit]
Description=Node Exporter

[Service]
Type=simple
DynamicUser=yes
ExecStart=/usr/local/bin/node_exporter

[Install]
WantedBy=multi-user.target